# python-scripts/sensor/service-example.conf
description "Sensor Stub"
author "CitOpt"

start on filesystem or runlevel [2345]
stop on runlevel [016]

chdir %dir%
exec python3 sensorstub.py -i %id% -t sensor/%id% -u %mqtturl%
respawn