# python-scripts/rmp-value-logger.service
[Unit]
Description=RMP value logger
Documentation=https://github.com/Rossojo/connde
After=mosquitto.service
Requires=mosquitto.service

[Service]
Type=simple
ExecStart=/opt/rmp/rmp-value-logger.sh
Restart=always

[Install]
WantedBy=multi-user.target